# overlay/usr/local/bin/mosquittoget
#!/bin/sh

# Get configuration setting from mosquitto configuration file

# shellcheck disable=SC1091,2154
. /etc/rc.subr && load_rc_config mosquitto

CONFIG="${mosquitto_config:-/usr/local/etc/mosquitto/mosquitto.conf}"
PWFILE="$(grep ^password_file "${CONFIG}" | awk '{print($2)}')"

case "${1}" in
  listener|port)
    grep ^listener "${CONFIG}" | awk '{print($2)}'
    ;;
  allow_anonymous)
    grep ^allow_anonymous "${CONFIG}" | awk '{print($2)}'
    ;;
  users)
    cut -d: -f1 "${PWFILE}"
    ;;
  *)
    echo "Unknown option: ${1}"
    ;;
esac
